# resource/cacheman.service
[Unit]
Description=Share pacman cache across hosts
Requires=avahi-daemon.service
Wants=network.target
After=avahi-daemon.service network.target

[Service]
Type=simple
ExecStart=/usr/bin/cacheman

[Install]
WantedBy=multi-user.target
